# serverbackup.service
[Unit]
Description=Runs serverbackup.py
Wants=serverbackup.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/serverbackup.py
User=root
